# supervisor.ini
[program:v-bot]
directory=/srv/v-bot
command=/srv/v-bot/v-bot
process_name=%(program_name)s
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/v-bot.log
stderr_logfile=/var/log/v-bot.log
