# supervisor.ini
[program:cuitclock]
directory=/srv/cuitclock
command=/srv/cuitclock/cuitclock
process_name=%(program_name)s
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cuitclock.log
stderr_logfile=/var/log/cuitclock.log
